# phtunnel/files/etc/hotplug.d/iface/30-oray-phtunnel
#!/bin/sh

[ "$INTERFACE" = "wan" ] && [ "$ACTION" = "ifup" ] && {
	logger -t phtunnel 'wan interface up, restarting phtunnel'
	/etc/init.d/phtunnel restart
}
