# phtunnel-apk/files/etc/hotplug.d/iface/30-oray-phtunnel
#!/bin/sh
# 30-oray-phtunnel
# 网络接口热插拔脚本

[ "$ACTION" = "ifup" ] || exit 0

# 检查 phtunnel 是否启用
if [ -f /etc/config/phtunnel ]; then
    enabled=$(uci get phtunnel.core.enabled 2>/dev/null || echo "0")
    
    if [ "$enabled" = "1" ]; then
        # 网络接口启动后，重启 phtunnel
        /etc/init.d/phtunnel restart
    fi
fi
